# debian/postinst
#!/bin/sh

set -e

. /etc/default/openmediavault
. /usr/share/openmediavault/scripts/helper-functions

LINK32="http://status.calibre-ebook.com/dist/linux32"
LINK64="http://status.calibre-ebook.com/dist/linux64"
FILE32="/tmp/calibre32.tar.bz2"
FILE64="/tmp/calibre64.tar.bz2"
INSTALL_DIR="/opt/calibre"

AMD64=false
ARCH="`/bin/uname -m`"

case "$ARCH" in
    *x86_64*)
        AMD64=true
        ;;
    *)
        AMD64=false
        ;;
esac

case "$1" in

    configure)

        if ! grep -q "^calibre:" /etc/passwd; then
            echo "Adding calibre user..."
            useradd --create-home --groups users --no-user-group --shell /usr/sbin/nologin calibre
        else
            echo "calibre user already exists."
        fi

        if ! omv_config_exists "/config/services/calibre"; then
            echo "Initial configuration."

            CALIBRE_HOME=$(getent passwd calibre | cut -d: -f6)

            object="<enable>0</enable>"
            object="${object}<mntentref></mntentref>"
            object="${object}<library-folder>${CALIBRE_HOME}</library-folder>"
            object="${object}<port>8080</port>"
            object="${object}<username></username>"
            object="${object}<password></password>"
            object="${object}<coversize></coversize>"
            object="${object}<showtab>0</showtab>"
            object="${object}<sharedfolderref></sharedfolderref>"

            omv_config_add_element "/config/services" "calibre" "${object}" true
        fi

        if [ -e /var/run/calibre ]; then
            echo "Removing old pid file..."
            rm -rf /var/run/calibre
        fi

        echo "Installing Calibre..."
        if [ ! -d "${INSTALL_DIR}" ]; then
            mkdir -p ${INSTALL_DIR}
            chown calibre:users ${INSTALL_DIR}
        fi

        if $AMD64; then
            if [ -f "${FILE64}" ]; then
                rm -f ${FILE64}
            fi
            wget ${LINK64} -O ${FILE64}

            echo "Extracting..."
            tar -xjof ${FILE64} -C ${INSTALL_DIR}
            rm -f ${FILE64}
        else
            if [ -f "${FILE32}" ]; then
                rm -f ${FILE32}
            fi
            wget ${LINK32} -O ${FILE32}

            echo "Extracting..."
            tar -xjof ${FILE32} -C ${INSTALL_DIR}
            rm -f ${FILE32}
        fi

        chown -R calibre:users ${INSTALL_DIR}

        insserv calibre

        if [ "$(omv_config_get "//services/calibre/enable")" != "1" ]; then
            update-rc.d calibre disable
        fi

        echo "Create config file..."
        /usr/share/openmediavault/mkconf/calibre

        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument" >&2
        exit 1
        ;;
esac

exit 0
